# client/restic-prune.timer
[Unit]
Description=Weekly Restic Key Pruning Timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
