# client/restic-backup.timer
[Unit]
Description=Daily Restic Backup Timer

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
